# cookbooks/vagrant_main/templates/default/magento.cron.erb
# /etc/cron.d/magento: crontab fragment for magento

*/5 * * * *  root  <%= @magentoroot %>/cron.sh