# cookbooks/vagrant_main/templates/default/magento.cron.erb
# /etc/cron.d/magento: crontab fragment for magento

* * * * *  root  "<%= @cron_sh %> <%= @cron_php %>"
